# root/etc/cron.monthly/UpdateArpScanVendors
#!/usr/bin/bash

cd /usr/share/arp-scan/
/usr/sbin/get-oui -u http://standards-oui.ieee.org/oui/oui.txt
/usr/sbin/get-iab -u http://standards-oui.ieee.org/iab/iab.txt
