# root/etc/cron.monthly/UpdateArpScanVendors
#!/usr/bin/bash
cd /usr/share/arp-scan/
get-oui -v 
get-iab -v
